# packages/sandbox/systemd/buildkitd.service
[Unit]
Description=BuildKit Daemon (cmux sandbox)
After=network.target dockerd.service
Requires=dockerd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/buildkitd
Restart=always
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
