# configs/systemd/cmux-worker.service
[Unit]
Description=Cmux Worker Service
After=cmux-ide.service
Requires=cmux-ide.service

[Service]
Type=simple
EnvironmentFile=-/opt/app/app.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=NODE_PATH=/builtins/build/node_modules
Environment=WORKER_PORT=39377
Environment=IS_SANDBOX=1
LogsDirectory=cmux
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/bin/bash -lc 'cd /root/workspace && exec /usr/bin/node /builtins/build/index.js'
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/cmux/worker.log
StandardError=append:/var/log/cmux/worker.log

[Install]
WantedBy=cmux.target
